# meta-edgeos/recipes-core/mender-expand/files/expand-partitions.service
[Unit]
Description=Expand EdgeOS partitions on first boot
DefaultDependencies=no
Before=sysinit.target
After=systemd-remount-fs.service
ConditionPathExists=!/var/lib/edgeos/partitions-expanded

[Service]
Type=oneshot
ExecStart=/usr/sbin/expand-partitions.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target