# meta-edgeos/recipes-core/edgeos-identity/files/edgeos-identity.service
[Unit]
Description=EdgeOS Device Identity and mDNS Setup
DefaultDependencies=no
Before=avahi-daemon.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/generate-uuid.sh
ExecStart=/usr/bin/update-mdns-uuid.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target