# meta-edgeos/recipes-support/expand-rootfs/files/expand-rootfs.service
[Unit]
Description=EdgeOS First Boot Rootfs Resize
After=nv.service
Before=nv-late-init.service
ConditionPathExists=!/etc/edgeos/.rootfs-resized

[Service]
Type=oneshot
# Create the trigger file that nv-late-init.sh checks
ExecStartPre=/bin/mkdir -p /etc/nv
ExecStartPre=/bin/touch /etc/nv/nvautoconfig
# Mark that we've triggered resize (one-time only)
ExecStart=/bin/mkdir -p /etc/edgeos
ExecStart=/bin/touch /etc/edgeos/.rootfs-resized
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
