# meta-edgeos/recipes-support/edgeos-services/files/90-edgeos-usbgadget.rules
ACTION=="bind", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_interface", DRIVER=="cdc_ether", RUN+="/usr/local/sbin/edgeos-usbgadget-resume.sh"
ACTION=="bind", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_interface", DRIVER=="ncm", RUN+="/usr/local/sbin/edgeos-usbgadget-resume.sh"
ACTION=="bind", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_interface", DRIVER=="usb_f_ncm", RUN+="/usr/local/sbin/edgeos-usbgadget-resume.sh"

# Handle wakeup from suspend
ACTION=="add", SUBSYSTEM=="usb", ATTR{power/wakeup}=="enabled", RUN+="/usr/local/sbin/edgeos-usbgadget-resume.sh"

# Handle resuming after system sleep
SUBSYSTEM=="power_supply", ACTION=="change", ATTR{status}=="Connected", RUN+="/usr/local/sbin/edgeos-usbgadget-resume.sh"

# Handle USB disconnect and reconnect
ACTION=="remove", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", RUN+="/bin/sleep 2 && /usr/local/sbin/edgeos-usbgadget-resume.sh"
ACTION=="add", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", RUN+="/bin/bash -c '/usr/local/sbin/edgeos-usbgadget-resume.sh >> /tmp/usb-gadget-debug.log 2>&1'"