# meta-edgeos/recipes-support/edgeos-services/files/edgeos-mdns-setup.service
[Unit]
Description=EdgeOS Device UUID and mDNS Setup
After=network.target
Before=avahi-daemon.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/generate-uuid.sh
ExecStart=/usr/local/sbin/update-edgeos-dns-service.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target 