# meta-edgeos/recipes-core/usb-gadget/files/usb-gadget-unbind.service

[Unit]
Description=Unbind USB Gadget
After=usb-gadget.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/usb-gadget.sh unbind
