# meta-edgeos/recipes-core/usb-gadget/files/usb-gadget-prepare.service

[Unit]
Description=Prepare USB Gadget (NCM-only)
After=sysinit.target local-fs.target systemd-udevd.service systemd-networkd.service
Wants=systemd-udevd.service systemd-networkd.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/usb-gadget.sh prepare

# Bind right away (autodetects the UDC), then bring usb0 admin-UP
ExecStartPost=/usr/sbin/usb-gadget.sh bind
ExecStartPost=/usr/libexec/usb0-force-up

# Clean teardown on shutdown/reboot
ExecStop=/usr/sbin/usb-gadget.sh destroy

RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
